# express-m3u-tuner.service
[Unit]
Description=Express M3U Tuner Service
After=network.target

[Service]
Type=simple
User=jellyfin
Group=jellyfin
WorkingDirectory=/opt/jellyfin-tuner
ExecStart=/usr/bin/node index.js
Restart=always

# Transcoding mode: none (direct copy), soft, qsv, nvenc, vaapi
Environment=TRANSCODE_MODE=none
Environment=ENABLE_EPG=true

# Supplementary groups for DVB and GPU access
SupplementaryGroups=video render

[Install]
WantedBy=multi-user.target
